# conserver.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    $Id$
-->
<service_bundle type='manifest' name='bx:conserver'>

	<service	name='bx/network/conserver'
		type='service'
		version='1'>

		<create_default_instance enabled='false' />

		<dependency	name='filesystem'
			grouping='require_all'
			restart_on='none'
			type='service'>
			<service_fmri value='svc:/system/filesystem/local' />
		</dependency>

		<dependency name='net-physical'
			grouping='require_all'
			restart_on='none'
			type='service'>
			<service_fmri value='svc:/network/physical' />
		</dependency>

		<exec_method name='start'
			type='method'
			exec='/opt/bx/heimdal/sbin/conserver.method start'
			timeout_seconds='300' />

		<exec_method name='stop'
			type='method'
			exec='/opt/bx/heimdal/sbin/conserver.method stop'
			timeout_seconds='60' />

		<property_group name='conserver' type='application'>
			<stability value='Unstable' />
			<propval name='options' type='astring' value='-d' />
		</property_group>

		<template>
			<common_name>
				<loctext xml:lang='C'>
				Conserver
				</loctext>
			</common_name>
		</template>
	</service>

</service_bundle>
